# rootfs-overlay/etc/systemd/system/input-event-daemon.service.d/override.conf
[Service]
StandardOutput=journal+console
ExecStart=
ExecStart=/usr/bin/input-event-daemon --no-daemon --verbose
